# deploy/pi-classroom-device.service
[Unit]
Description=Raspberry Pi Classroom Device MVP
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pi-classroom-device
ExecStart=/usr/bin/python3 -m device_app --config /opt/pi-classroom-device/config/device.sample.json --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
